# debian/platform-modules-phalanx.postinst
depmod -a

# Enable bmc vlan_ip service
systemctl enable bmc_vlan.service

systemctl enable platform-modules-phalanx.service
systemctl start platform-modules-phalanx.service

# Enable platform-sync_bmc-timer
# systemctl enable sync_bmc.timer
# systemctl start sync_bmc.timer
